# src/rc.d/abuseipdbchecker
#!/bin/sh
#
# PROVIDE: abuseipdbchecker
# REQUIRE: NETWORKING syslogd
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable abuseipdbchecker:
# abuseipdbchecker_enable="YES"
#

. /etc/rc.subr

name="abuseipdbchecker"
rcvar="abuseipdbchecker_enable"

load_rc_config $name

: ${abuseipdbchecker_enable="NO"}
: ${abuseipdbchecker_interval="3600"}

pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
command_args="-f -P ${pidfile} -t ${name} /usr/local/opnsense/scripts/AbuseIPDBChecker/checker.py"
abuseipdbchecker_cron="/usr/local/etc/cron.d/abuseipdbchecker"

abuseipdbchecker_prestart()
{
    # Generate cron job
    echo "# DO NOT EDIT THIS FILE - Generated by rc.d/${name}" > ${abuseipdbchecker_cron}
    echo "0 * * * * root /usr/local/opnsense/scripts/AbuseIPDBChecker/checker.py >> /var/log/abuseipdb_checker.log 2>&1" >> ${abuseipdbchecker_cron}
    
    # Ensure the script exists and is executable
    if [ ! -x /usr/local/opnsense/scripts/AbuseIPDBChecker/checker.py ]; then
        echo "Error: Checker script not found or not executable"
        return 1
    fi
    
    # Ensure log file exists and has proper permissions
    touch /var/log/abuseipdb_checker.log
    chmod 640 /var/log/abuseipdb_checker.log
    
    return 0
}

abuseipdbchecker_poststop()
{
    # Remove cron job
    rm -f ${abuseipdbchecker_cron}
    return 0
}

start_cmd="${name}_start"
stop_cmd="${name}_stop"

abuseipdbchecker_start()
{
    abuseipdbchecker_prestart || return 1
    
    # Initial run
    /usr/local/opnsense/scripts/AbuseIPDBChecker/checker.py >> /var/log/abuseipdb_checker.log 2>&1
    
    echo "Started ${name} service (running on cron schedule)"
}

abuseipdbchecker_stop()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        kill -TERM "${pid}" 2>/dev/null && echo "Stopped ${name} process"
        rm -f "${pidfile}"
    fi
    
    abuseipdbchecker_poststop
    echo "Stopped ${name} service"
}

run_rc_command "$1"
